# novel_coronavirus/conf/spider.conf
[supervisord]
nodaemon=true


[program:2019ncov]
user=root
directory=/opt/2019ncov/novel_coronavirus
command=/usr/local/bin/scrapyd --logfile=/tmp/spider.log --pidfile=spider.pid
