# docker/conf/supervisor.conf
[program:2019ncov]
environment =
    DJANGO_SETTINGS_MODULE=ncov.settings.dev,
user=root
directory=/opt/2019ncov/ncov
command=/usr/local/bin/gunicorn -c /etc/default/gunicorn ncov.wsgi